# .config/systemd/user/meeting.service
[Unit]
Description=Meeting time notification

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send -u critical "meeting"
RemainAfterExit=no
