# .config/systemd/user/meeting.timer
[Unit]
Description=Meeting time

[Timer]
OnCalendar=Mon..Fri 13:55

[Install]
WantedBy=timers.target
