# .config/systemd/user/show_time.timer
[Unit]
Description=show time notification

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
