# .config/systemd/user/show_time.service
[Unit]
Description=show time notification

[Service]
Type=oneshot
ExecStart=%h/.bin/misc/show_time.sh
RemainAfterExit=no
